# systemd/mad-daemon.service
[Unit]
Description=mad privileged daemon
After=network.target

[Service]
ExecStart=/usr/bin/mad daemon
Restart=on-failure
RestartSec=2

# Runtime dirs
RuntimeDirectory=mad
RuntimeDirectoryMode=0750
StateDirectory=mad
StateDirectoryMode=0750

# Privileges
User=root
Group=mad
AmbientCapabilities=CAP_NET_ADMIN CAP_DAC_OVERRIDE CAP_CHOWN CAP_FOWNER CAP_SETUID CAP_SETGID

# Hardening
# The daemon shells out to useradd / usermod / chpasswd / passwd to
# manage Linux accounts during enrollment, and writes to /etc/ssh/mad_krl
# on every cert revocation. Those tools touch /etc/passwd, /etc/group,
# /etc/shadow, /etc/gshadow, /etc/subuid, /etc/subgid, /home/*,
# /var/spool/mail/*, /etc/ssh/. ProtectSystem=strict and
# ProtectHome=read-only would block all of that; leave them off and rely
# on the explicit AmbientCapabilities + User=root model instead.
NoNewPrivileges=no
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
